# systemd/amanuense.path
[Unit]
Description=Recarrega Whisper Dictation ao alterar config.toml
Documentation=https://github.com/renatoobarros/amanuense

[Path]
PathChanged=%h/.config/amanuense/config.toml
PathModified=%h/.config/amanuense/config.toml
Unit=amanuense-restart.service

[Install]
WantedBy=default.target
