# systemd/amanuense-restart.service
[Unit]
Description=Reinício do Whisper Dictation após mudança de configuração
Documentation=https://github.com/renatoobarros/amanuense

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --user try-restart amanuense.service
