# systemd/amanuense.service
[Unit]
Description=Amanuense Daemon
Documentation=https://github.com/renatoobarros/amanuense
After=graphical-session.target pipewire.service pipewire-pulse.service
Requires=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/amanuense daemon
Restart=on-failure
RestartSec=5

Environment=RUST_LOG=warn

NoNewPrivileges=true

[Install]
WantedBy=graphical-session.target
