# systemd/whisper-dictate-restart.service
[Unit]
Description=Reinício do Whisper Dictation após mudança de configuração
Documentation=https://github.com/renatoobarros/whisper-dictate

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --user try-restart whisper-dictate.service
